# src/frontend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:php-fpm]
command=php-fpm7.4 -F

[program:nginx]
command=nginx -g 'daemon off;'
